# 1180109-daniel-azevedo/p2.1/configuration_files/ansible-playbook/playbooks/config_files/todd/unit/todd.service
[Unit]
Description=My TODD cogsi app
After=network.target

[Service]
Type=simple
User=nagios
WorkingDirectory=/home/dan/todd-cogsi
ExecStart=/usr/local/bin/start_todd

[Install]
WantedBy=multi-user.target